# pwnarc.service
[Unit]
Description=arc passwod manager

[Service]
Type=simple
ExecStart=/home/pi/scripts/arc.sh
Restart=always

[Install]
WantedBy=multi-user.target
