# pwnshutdown.service
[Unit]
Description=pwnshutdown service 

[Service]
Type=oneshot
ExecStart=/home/pi/scripts/pwnshutdown.sh