# pwntimecheck.service
[Unit]
Description=pwntime check
Type=oneshot

[Service]
ExecStart=/home/pi/scripts/pwntimecheck.sh
